# docs/systemd/linear-org-sync.service
[Unit]
Description=Sync Linear issues to org-mode files

[Service]
Type=oneshot
WorkingDirectory=%h/src/github/tellett/monorepo
ExecStart=/home/tellett/go/bin/bazel run //py/linear_org_sync -- --config %h/org/linear/config.yaml
StandardOutput=journal
StandardError=journal
