# docs/systemd/linear-org-sync.timer
[Unit]
Description=Run Linear org sync periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
